# appuser.cron
*/5 * * * *     /home/appuser/slashdev/bin/slashdev-update.sh
