# appuser.cron
*/5 * * * *     ~/slashdev/bin/slashdev-update.sh
